# chat-server/bin/chat-server.service
<<<<<<< HEAD
; /etc/systemd/system/chat-server.service
=======
# /etc/systemd/system/chat-server.service

>>>>>>> bb6f6a7694c61b4085a7f08cf8c6e00e55b3bbdf
[Unit]
Description=Node.js chat-server

[Service]
PIDFile=/tmp/chat-server-99.pid
User=<user>
Group=<group>
Restart=always
KillSignal=SIGTERM
WorkingDirectory=/martijncasteel.nl/chat-server/
ExecStart=/martijncasteel.nl/chat-server/bin/chat-server.js

[Install]
WantedBy=multi-user.target
