# chat-server/bin/chat-server.service
; /etc/systemd/system/chat-server.service
[Unit]
Description=Node.js chat-server

[Service]
PIDFile=/tmp/chat-server-99.pid
User=<user>
Group=<group>
Restart=always
KillSignal=SIGTERM
WorkingDirectory=/martijncasteel.nl/chat-server/
ExecStart=/martijncasteel.nl/chat-server/bin/chat-server.js

[Install]
WantedBy=multi-user.target
